# service/growtent-web.service
[Unit]
Description=Grow Tent Web UI (Flask/Gunicorn)
After=network-online.target
Wants=growtent.service

[Service]
Type=simple
User=matthias
WorkingDirectory=/home/matthias/grow_tent_controller/flask-app
Environment=PYTHONUNBUFFERED=1
# app module = app:app within flask-app/app.py
ExecStart=/home/matthias/grow_tent_controller/.grow_tent_web/bin/gunicorn -b 0.0.0.0:8000 app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
